# apps/s6/service/emacs/run
#!/usr/bin/env execlineb -P

fdmove -c 2 1
exec fish -l -c "emacs --daemon --init-directory /Users/astephe9/.config/emacs/dotemacs"
