# emacs.service
[Unit]
Description=Emacs text editor
Documentation=info:emacs man:emacs(1) https://gnu.org/software/emacs/

[Service]
Type=simple
ExecStart=/home/adam/.nix-profile/bin/emacs --fg-daemon
ExecStop=/home/adam/.nix-profile/bin/emacsclient --eval "(kill-emacs)"
Environment=SSH_AUTH_SOCK=%t/keyring/ssh
Environment=WAYLAND_DISPLAY=wayland-1
Restart=on-failure

[Install]
WantedBy=default.target
