# kanshi.service
[Unit]
Description=kanshi
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/kanshi

[Install]
WantedBy=river-session.target
