# kanshi.service
[Unit]
Description=kanshi

[Service]
ExecStart=/usr/bin/kanshi

[Install]
WantedBy=default.target
