# apps/nix/non-nixos/nix-daemon.openrc
#!/sbin/openrc-run

name="nix-daemon"
description="Nix Daemon"
supervisor="supervise-daemon"
command="/nix/var/nix/profiles/default/bin/nix-daemon"
pidfile="/run/${RC_SVCNAME}.pid"
command_background=true
